# support/systemd/piwx.service.example
[Unit]
Description=METAR Display
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/local/piwx/bin/piwx -s
Restart=on-abort

[Install]
WantedBy=multi-user.target
